# services/enes100ml.service
[Unit]
Description=Enes100ML startup sequence
After=network.target

[Service]
ExecStart=/usr/bin/mlboot.sh

[Install]
WantedBy=default.target
